# .config/systemd/user/set-wallpaper.service
[Unit]
Description=Set GNOME wallpaper based on time

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/set-wallpaper.sh
